# systemd/aqm-nest.service
[Unit]
Description=AQM Nest audio alerts
After=network-online.target aqm-web.service
Wants=network-online.target

[Service]
Type=simple
# Adjust User and the paths to match your Pi.
User=pi
WorkingDirectory=/home/pi/air-quality-monitor
ExecStart=/home/pi/air-quality-monitor/.venv/bin/python /home/pi/air-quality-monitor/scripts/nest_alert.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
